# custom-build/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgresql]
command=/opt/dartfx/scripts/runtime/start-postgres.sh
user=postgres
autostart=true
autorestart=true
startsecs=5
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.err.log
stopsignal=TERM

[program:ollama]
command=/usr/local/bin/ollama serve
user=dartfx
directory=/home/dartfx
environment=HOME="/home/dartfx",USER="dartfx",OLLAMA_MODELS="/var/lib/ollama/models"
autostart=true
autorestart=true
startsecs=5
priority=20
stdout_logfile=/var/log/supervisor/ollama.log
stderr_logfile=/var/log/supervisor/ollama.err.log
stopsignal=TERM
stopasgroup=true
killasgroup=true
